# rootfs/etc/services.d/nginx/run
#!/bin/bash
mkdir -p /tmp/nginx

exec /usr/sbin/nginx
